# filament-queue-worker.service
[Unit]
Description=Laravel Queue Worker for Filament App
After=network.target mysql.service redis.service
Requires=mysql.service redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/filament-app
ExecStart=/usr/bin/php /var/www/filament-app/artisan queue:work redis --queue=polling,scheduling,default --sleep=3 --tries=3 --max-time=3600 --memory=512
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=LARAVEL_ENV=production

# Resource limits
LimitNOFILE=65536
MemoryMax=1G

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/filament-app/storage

[Install]
WantedBy=multi-user.target