# filament-polling-monitor.service
[Unit]
Description=Filament Polling System Monitor
After=network.target mysql.service redis.service filament-queue-worker.service
Requires=mysql.service redis.service

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/filament-app
ExecStart=/usr/bin/php /var/www/filament-app/artisan polling:reliable start
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=LARAVEL_ENV=production